# cron.d/scraper
# Each Monday at 08:00
0 8 * * 1   root    cd /opt/ef_schedule; /usr/bin/node scraper.js
